# assessment item 1/catkin_ws/devel/etc/catkin/profile.d/90-uol_cmp9767m_base.sh
# extend the Gazebo model path with our models
export GAZEBO_MODEL_PATH=/home/computing/catkin_ws/src/CMP9767M/uol_cmp9767m_base/models:/home/computing/catkin_ws/install/share/uol_cmp9767m_base/models:$GAZEBO_MODEL_PATH
